# var/spool/cron/crontabs/root
*/30 * * * * bash /home/nodo/update-all.sh
0 * * * * bash /home/nodo/execScripts/monero-wallet-rpc-sweep.sh
0 1 * * * bash /home/nodo/update-banlists.sh
