# var/spool/cron/crontabs/root
*/30 * * * * bash /home/nodo/update-all.sh
