# etc/systemd/system/webui.service
[Unit]
Description=Monero Web Interface
After=network.target

[Service]
ExecStart=/home/nodo/webui/venv/bin/python /home/nodo/webui/webapp.py 127.0.0.1 8082
WorkingDirectory=/home/nodo/webui/
Type=exec
Restart=always
RestartSec=5
User=nodo
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
